# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/bin/bash -c "/app/ssh.sh"

[program:python]
command=/bin/bash -c "/app/train-py.sh"

[program:tensorboard]
command=/bin/bash -c "/app/start-tb.sh"