# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:python]
command=/bin/bash -c "/app/train-py.sh"
